# systemd/webshell.service
[Unit]
Description=Web Shell
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/project/web_worker
EnvironmentFile=/etc/webshell.env
ExecStartPre=/usr/local/bin/openresty -p /opt/project/web_worker/openresty/ -c conf/nginx.conf -t
ExecStart=/bin/sh -c 'set -e; if [ -s /opt/project/web_worker/openresty/logs/openresty.pid ] && kill -0 "$(cat /opt/project/web_worker/openresty/logs/openresty.pid)" 2>/dev/null; then /usr/local/bin/openresty -p /opt/project/web_worker/openresty/ -c conf/nginx.conf -s reload; else /usr/local/bin/openresty -p /opt/project/web_worker/openresty/ -c conf/nginx.conf; fi; exec /opt/project/web_worker/webshell'
ExecStopPost=-/bin/sh -c '/usr/local/bin/openresty -p /opt/project/web_worker/openresty/ -c conf/nginx.conf -s stop; for i in 1 2 3 4 5; do [ ! -s /opt/project/web_worker/openresty/logs/openresty.pid ] || ! kill -0 "$(cat /opt/project/web_worker/openresty/logs/openresty.pid)" 2>/dev/null && exit 0; sleep 0.2; done'
KillMode=process
Restart=always
RestartSec=2

[Install]
WantedBy=multi-user.target
